# etc/supervisor/rootpassword.conf
[program:rootpassword]
directory=/usr/local/bin/
command=/usr/local/bin/set-root-password.sh
autostart=true
autorestart=false
redirect_stderr=true