# src/main/docker/5.6.28-apache/etc/supervisor/apache.conf
[program:apache]
directory=/usr/local/bin/
command=/usr/local/bin/apache2-foreground
autostart=true
autorestart=true
redirect_stderr=true